# parts/languages/php/files/service/log/run
#!/bin/sh

if [ ! -d "/log/php/service" ]; then
  mkdir -p /log/php/service
  chown -R app: /log/php/service
fi

exec multilog t /log/php/service

